# scripts/amd-tunnel/volc/amd-agent-tunnel.service
[Unit]
Description=RoboCloud AMD agent tunnel (control+data plane, volc -> AMD)
Documentation=file:///docs/RobotCloud_AMD_Tunnel_Deployment.md
After=network-online.target
Wants=network-online.target

[Service]
# 以专用非 root 账号运行；凭据全部在 /var/lib/rc-tunnel/.ssh 下
User=rc-tunnel
Group=rc-tunnel
Environment=AUTOSSH_GATETIME=0
# -M 0 关闭 autossh 自带监控端口，改用 ssh 的 ServerAlive 保活（见 ssh config）
# -F 指定 rc-tunnel 的 ssh config；-L 把 volc:127.0.0.1:__AGENT_PORT__ 转发到 AMD 内部同端口
ExecStart=/usr/bin/autossh -M 0 -N -F /var/lib/rc-tunnel/.ssh/config \
  -o ExitOnForwardFailure=yes \
  -L 127.0.0.1:__AGENT_PORT__:localhost:__AGENT_PORT__ \
  amd-gpu
Restart=always
RestartSec=5

# --- 加固 ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/rc-tunnel
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
